# svc/manifest/minecraft.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="minecraft">
  <service name="application/minecraft" type="service" version="1">
    <create_default_instance enabled="false" />
    <single_instance />
    <dependency name="network" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/milestone/network:default" />
    </dependency>
    <dependency name="filesystem" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/system/filesystem/local" />
    </dependency>
    <method_context working_directory="/opt/minecraft/server">
      <method_credential user="minecraft" group="minecraft" />
      <method_environment>
        <envvar name="PATH" value="/opt/local/bin:/opt/local/sbin:/usr/bin:/usr/sbin:/usr/local/bin" />
      </method_environment>
    </method_context>
    <exec_method type="method" name="start" exec="/opt/minecraft/svc/method/minecraft start" timeout_seconds="60" />
    <exec_method type="method" name="stop" exec="/opt/minecraft/svc/method/minecraft stop" timeout_seconds="60" />
    <exec_method type="method" name="refresh" exec="/opt/minecraft/svc/method/minecraft refresh" timeout_seconds="60" />
    <property_group name="startd" type="framework">
      <propval name="duration" type="astring" value="contract" />
    </property_group>
    <property_group name="minecraft" type="application">
      <propval name="java_opts" type="astring" value="-d64 -Xincgc -Djava.security.egd=/dev/./urandom -XX:-UseVMInterruptibleIO -XX:+UseConcMarkSweepGC -XX:+CMSIncrementalPacing -XX:ParallelGCThreads=4 -XX:+AggressiveOpts" />
      <propval name="mc_flavor" type="astring" value="minecraft" />
      <propval name="min_size" type="astring" value="512M" />
      <propval name="max_size" type="astring" value="1536M" />
    </property_group>
    <stability value="Evolving" />
    <template>
      <common_name>
        <loctext xml:lang="C">Minecraft service</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
